# backend/ally-cpu-boost-disabler-cap-refresh.service
[Unit]
Description=Refresh amd-pstate CPU cap after power source change (ally-cpu-boost-disabler)

[Service]
Type=oneshot
# Bazzite SELinux blocks executing scripts in /var/lib directly; bash reads the file instead.
ExecStart=/usr/bin/bash /var/lib/ally-cpu-boost-disabler/refresh-cpu-cap.sh
ExecStart=/usr/bin/sleep 15
ExecStart=/usr/bin/bash /var/lib/ally-cpu-boost-disabler/refresh-cpu-cap.sh
